# systemd/dirigera-computed-schedule.service
[Unit]
Description=Apply computed DIRIGERA Adaptive schedule
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=piadmin
Group=piadmin
WorkingDirectory=/opt/dirigera-readaptive
Environment=PYTHONPATH=/opt/dirigera-readaptive/src
ExecStart=/opt/dirigera-readaptive/.venv/bin/python -m dirigera_readaptive.schedule_profiles_cli --config /opt/dirigera-readaptive/config.yaml
